# crontab
10 0 * * * /app/backup_script.sh >> /var/log/cron.log 2>&1
